# home/madhur/.config/systemd/user/everymonth.service.d/notifications.conf
[Unit]
OnSuccess=notify@%N.service
OnFailure=notify-fail@%N.service
