# home/madhur/.config/systemd/user/postgres_exporter.service.d/notifications.conf
[Unit]
OnSuccess=notify@%N.service
OnFailure=notify-fail@%n.service
