# home/madhur/.config/systemd/user/authelia-audit.timer
[Unit]
Description=Run Authelia login audit daily at 08:00

[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
