# home/madhur/.config/systemd/user/desktop-screenshot.service
[Unit]
Description=Desktop Screenshot Service
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/take-screenshot.sh

