# home/madhur/.config/systemd/user/everyweek.service
[Unit]
Description=Script to perform weekly tasks

[Service]
Type=oneshot
ExecStart=/home/madhur/scripts/every_week.sh
ExecStartPre=/bin/sleep 120