# home/madhur/.config/systemd/user/everyhalfhour.service
[Unit]
Description=Script to do anything every half hour

[Service]
Type=oneshot
ExecStart=/home/madhur/scripts/every_half_hour.sh