# home/madhur/.config/systemd/user/everyday.service
[Unit]
Description=Script to do anything every day

[Service]
Type=oneshot
ExecStart=/home/madhur/scripts/every_day.sh
ExecStartPre=/bin/sleep 120