# home/madhur/.config/systemd/user/everyday.timer
[Unit]
Description=Run unit every day

[Timer]
OnBootSec=15min
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
