# home/madhur/.config/systemd/user/everyhour.service
[Unit]
Description=Script which runs every hour

[Service]
Type=oneshot
ExecStart=/home/madhur/scripts/wallpaper.sh
