# home/madhur/.config/systemd/user/everymonth.service
[Unit]
Description=Script to perform monthly tasks

[Service]
Type=oneshot
ExecStart=/home/madhur/scripts/every_month.sh
