# home/madhur/.config/systemd/user/everymonth.timer
[Unit]
Description=Run unit every month

[Timer]
Persistent=true
OnCalendar=monthly
RandomizedDelaySec=3600


[Install]
WantedBy=timers.target