# home/madhur/.config/systemd/user/everyhour.timer
[Unit]
Description=Run unit every hour
Requires=everyhour.service

[Timer]
OnBootSec=1min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
