# home/madhur/.config/systemd/user/everyfiveminutes.service
[Unit]
Description=Take screenshot every five minutes

[Service]
Type=oneshot
ExecStart=/home/madhur/scripts/take-screenshot.sh --timer
