# home/madhur/.config/systemd/user/everyweek.timer
[Unit]
Description=Run unit every week

[Timer]
Persistent=true
OnCalendar=weekly
RandomizedDelaySec=3600


[Install]
WantedBy=timers.target