# home/madhur/.config/systemd/user/everyhour.service
[Unit]
Description=Script to do anything every hour

[Service]
Type=oneshot
ExecStart=/home/madhur/scripts/every_hour.sh
