# home/madhur/.config/systemd/user/everyhour.timer
[Unit]
Description=Run unit every hour

[Timer]
#OnBootSec=1min
OnCalendar=hourly
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
