# home/madhur/.config/systemd/user/idle-shutdown.timer
[Unit]
Description=Timer to check system idle time every 15 minutes

[Timer]
# Run every 15 minutes
OnCalendar=*:0/15
# Run immediately if timer was missed (e.g., system was off)
Persistent=false
# Add a small random delay to avoid thundering herd (0-60 seconds)
RandomizedDelaySec=60

[Install]
WantedBy=default.target

