# home/madhur/.config/systemd/user/idle-shutdown.timer
[Unit]
Description=Timer to check system idle time every minute

[Timer]
# Run every minute (also drives the zjstatus idle-countdown bar push)
OnCalendar=*:*:00
# Fire crisply on the minute boundary
AccuracySec=1s
# No random delay: we want a predictable once-per-minute push
RandomizedDelaySec=0
# Run immediately if timer was missed (e.g., system was off)
Persistent=false

[Install]
WantedBy=default.target
