# home/madhur/.config/systemd/user/update-dns.timer
[Unit]
Description=Run Update DNS service every 5 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
