# home/madhur/.config/systemd/user/xvfb.service
[Unit]
Description=Xvfb virtual framebuffer on display :98

[Service]
Type=simple
ExecStart=/usr/bin/Xvfb :98 -screen 0 1280x720x24 -ac
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
