# home/madhur/.config/systemd/user/zellij-bar.timer
[Unit]
Description=Refresh zellij zjstatus bar (load + clock) every minute

[Timer]
# Fire on the minute boundary so the clock reads HH:MM accurately.
OnCalendar=*:*:00
AccuracySec=1s
RandomizedDelaySec=0
Persistent=false

[Install]
WantedBy=default.target
